# .config/systemd/user/backuper-restic.timer
[Unit]
Description=Run backup script

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h
 
[Install]
WantedBy=timers.target
