# .config/systemd/user/backuper-restic.service
[Unit]
Description=restic backup

[Service]
ExecStart=/home/amish/bin/run-backup.sh
TimeoutStopSec=30
