# .config/systemd/user/plasma-i3.service
[Unit]
Description=Launch Plasma with i3
After=plasma-workspace.target

[Service]
ExecStart=/usr/bin/i3
Restart=on-failure

[Install]
WantedBy=plasma-workspace.target
