# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Name=xbindkeys
Terminal=false
Type=Application
Exec=xbindkeys

